# server/docker/etc/services.d/boardr/run
#!/usr/bin/with-contenv sh
sudo -g boardr -u boardr /usr/src/app/bin/boardr start